# ezrec-backend.service
[Unit]
Description=EZREC Backend Service - Soccer Recording System
Documentation=https://github.com/naeimsalib/EZREC-BackEnd
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=michomanoly14892
Group=michomanoly14892
WorkingDirectory=/opt/ezrec-backend
ExecStartPre=/bin/bash -c 'echo "ðŸ›¡ï¸ Protecting camera for EZREC..."'
ExecStartPre=/bin/bash -c 'sudo fuser -k /dev/video0 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'sudo fuser -k /dev/video1 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'sudo fuser -k /dev/video2 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'echo "âœ… Camera protection active"'
ExecStart=/opt/ezrec-backend/venv/bin/python3 main.py
ExecReload=/bin/kill -USR1 $MAINPID

# Environment
Environment=PYTHONPATH=/opt/ezrec-backend
EnvironmentFile=-/opt/ezrec-backend/.env

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ezrec-backend
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
MemoryMax=1G
CPUQuota=80%

# Restart policy
Restart=always
RestartSec=10
TimeoutStartSec=60
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ezrec-backend

# Additional video group access for camera
SupplementaryGroups=video

[Install]
WantedBy=multi-user.target 